# cellapp.service
[Unit]
Description=CellApp Flask Application
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/var/www/cellapp
Environment="PATH=/var/www/cellapp/venv/bin"
EnvironmentFile=/var/www/cellapp/.env
ExecStart=/var/www/cellapp/venv/bin/gunicorn \
    --workers 4 \
    --worker-class sync \
    --bind 127.0.0.1:5001 \
    --timeout 60 \
    --access-logfile /var/log/cellapp/access.log \
    --error-logfile /var/log/cellapp/error.log \
    --log-level info \
    wsgi:application

# Restart policy
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
